# package/common/usr/lib/systemd/system/simbiota-update.timer
[Unit]
Description=Simbiota database updater
After=network.targe

[Timer]
OnCalendar=0/1:00:00

[Install]
WantedBy=timers.target
